# graba-video-1080.service

[Unit]
Description=video-1080
After=multi-user.target

[Service]
TimeoutStartSec=0
ExecStart=/home/pi/comtrolar--navio-video-telemetria-master/video-1080
Restart=always

[Install]
WantedBy=multi-user.target

